# systemd/wireplumber-restart.service
[Unit]
Description=Fix audio profile on boot (Intel SOF boot race)
After=wireplumber.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c '\
    sleep 2 && \
    systemctl --user restart wireplumber && \
    for i in $(seq 1 15); do \
        card=$(pactl list cards short 2>/dev/null | awk "/skl_hda_dsp_generic/{print \$2}"); \
        [ -n "$card" ] && pactl set-card-profile "$card" "HiFi (HDMI1, HDMI2, HDMI3, Mic1, Mic2, Speaker)" && exit 0; \
        sleep 1; \
    done; \
    exit 1'

[Install]
WantedBy=default.target
